# common/scripts/kvm/kvm-lava.conf

start on runlevel [23]
stop on runlevel [!23]
console output
script
    echo "kvm-boot-1: pass"
    echo "Test hackbench on guest"
    sh /root/test-rt-tests.sh
    poweroff
end script
