# common/scripts/kvm/kvm-lava.conf

start on runlevel [23]
stop on runlevel [!23]
console output
script
    sh /usr/bin/test-guest.sh
    poweroff
end script
